# metro-map.service
[Unit]
Description=DC Metro Map LED Display
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/metro-map
Environment=PYTHONUNBUFFERED=1
Environment=TMPDIR=/run/metro-map

# Create RAM-based directory for runtime data
RuntimeDirectory=metro-map
RuntimeDirectoryMode=0755

# Use the virtualenv Python
ExecStart=/home/pi/metro-map/.venv/bin/python app.py

# Reduce writes to SD card
ReadWritePaths=/run/metro-map
ProtectSystem=strict
ProtectHome=read-only

# Restart settings
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
