# scripts/metro-map-restart.timer
[Unit]
Description=Daily restart of DC Metro Map LED Display
RefuseManualStart=no
RefuseManualStop=no

[Timer]
# Run daily at 4 AM (when metro is closed)
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=300
Unit=metro-map.service

[Install]
WantedBy=timers.target
